# docker/crontab
# Quant Hub cron — weekday breakout scan at 5:17 PM ET
17 17 * * 1-5 root . /etc/environment; quant-daily --universe sp500 >> /app/logs/cron.log 2>&1

# Weekly swing scan — Friday 6:17 PM ET (10y / 1wk OHLCV + email)
17 18 * * 5 root . /etc/environment; quant-swing --universe sp500 >> /app/logs/cron.log 2>&1
